# rootfs/etc/services.d/rtlamr2mqtt/run

# wait for rtl_tcp
bashio::net.wait_for 1234

bashio::log.info "Starting rtlamr...."